# docker/zookeeper/zookeeper.conf
[program:zookeeper_pre]
command=/opt/zookeeper/bin/pre.sh
process_name=%(program_name)s
autostart=true
autorestart=true
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0
priority=100

[program:zookeeper]
command=/opt/zookeeper/bin/zkServer.sh start-foreground
environment=ZOO_LOG_DIR=/var/log/zookeeper,ZOO_LOG4J_PROP='INFO,ROLLINGFILE'
process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=TERM
log_level=debug
user=root
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0
priority=200
