# docker/backup/rsync.conf
[program:rsync]
command = /usr/bin/rsync --no-detach --daemon --config="/etc/rsyncd.conf"
redirect_stderr = true
autorestart=true
